# bergwerk-wiki/crontab
# crontab

1 1 * * * bash -c "cd /var/www/html/w; php maintenance/run.php deleteArchivedRevisions --delete"
